# docker/opendkim/rootfs/etc/s6-overlay/s6-rc.d/createOpenDkimConfig/up
foreground { echo "Creating OpenDkim configuration by calling /usr/local/bin/createOpenDkimConfig.sh" }
/usr/local/bin/createOpenDkimConfig.sh